# src/wifi_reconnect.cron
# Run the wifi_reconnect script every minute
* *   * * *   root    /usr/local/bin/wifi_reconnect.sh
